# etc/systemd/system/android-mount.service.d/10-nogstdroid.conf
[Service]
ExecStartPre=/usr/bin/mount --bind /dev/null /usr/lib/aarch64-linux-gnu/gstreamer-1.0/libgstdroid.so
