# usr/lib/systemd/user/droidian-sysfs-torch.service
[Unit]
Description=Sysfs flashlight controller
After=org.droidian.Flashlightd.service

[Service]
Type=simple
ExecStart=/opt/droidian-sysfs-torch/bin/droidian-sysfs-torch
Restart=on-failure
